# deploy/auto-run.service
[Unit]
Description=Auto-Run: 中韩出海数智港 Backend API
Documentation=https://github.com/eagle13579/china-korea-digital-port
After=network.target

[Service]
Type=simple
User=opc
Group=opc
WorkingDirectory=/var/www/china-korea-digital-port
EnvironmentFile=/var/www/china-korea-digital-port/.env
Environment=PATH=/var/www/china-korea-digital-port/venv/bin:/usr/local/bin:/usr/bin:/bin
ExecStart=/var/www/china-korea-digital-port/venv/bin/uvicorn backend.main:app --host 127.0.0.1 --port 8000 --workers 1 --proxy-headers --forwarded-allow-ips="127.0.0.1"
Restart=always
RestartSec=5
StartLimitInterval=0
StartLimitBurst=5
StandardOutput=append:/var/log/ckdp/backend.log
StandardError=append:/var/log/ckdp/backend-error.log
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=off
ProtectHome=false
ReadWritePaths=/var/www/china-korea-digital-port /var/log/ckdp

[Install]
WantedBy=multi-user.target
